# falkor_bringup/upstart/robot.conf
description	"bringup robot"

start on runlevel [2345]

console output
respawn

exec robot-start
pre-stop exec robot-stop
