# falkor_bringup/upstart/robot.conf
description	"bringup robot"

start on net-device-up IFACE=wlan0
stop on net-device-down IFACE=wlan0

console output
respawn

exec robot-start
pre-stop exec robot-stop